# supervisor.conf
# ============================================================
# Supervisor 配置 — 苏州赛冠后端
# 文件路径:/etc/supervisor/conf.d/topchampion.conf
# 启用:sudo supervisorctl reread && sudo supervisorctl update
# ============================================================

[program:topchampion-backend]
command=/opt/topchampion/backend/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 2
directory=/opt/topchampion/backend
user=www-data
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10

stdout_logfile=/var/log/topchampion-backend.out.log
stderr_logfile=/var/log/topchampion-backend.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

environment=PYTHONUNBUFFERED=1
